# arm/configs/rc.local
#!/bin/sh
tmux new-session -d
tmux neww /home/sources/kvmtests/slave.py
